# templates/rocket.service.erb
<%- # guards
authstring = @authsource.to_s != '' ? "?authSource=#{@authsource}" : nil
replstring = @mongo_replset.to_s != '' ? (authstring ? '&' : '?') + "replicaSet=#{@mongo_replset}" : nil
-%>
[Unit]
Description=Rocket.Chat service

[Service]
Type=forking
<%#- Wait for the forked process to listen on a socket and make sure it is still alive -%>
ExecStart=/bin/sh -e -c '/usr/local/bin/node --max_old_space_size=<%= @max_old_space_size %> <%= @destination %>/bundle/main.js & while ! lsof -iTCP -sTCP:LISTEN -a -p $!; do kill -0 $!; sleep 1; done'
Restart=always
# Restart service after 10 seconds if node service crashes
RestartSec=10
LogsDirectory=rocketchat
StandardOutput=append:/var/log/rocketchat/out.log
StandardError=append:/var/log/rocketchat/err.log
SyslogIdentifier=rocket.chat
Environment=INSTANCE_IP=<%= @instance_ip || @ipaddress %>
Environment=PORT=<%= @instance_count == 1 ? @port : '%i' %>
Environment=MONGO_URL=mongodb+srv://<%= @mongo_host %>/<%= @database_name %><%= authstring %><%= replstring %>&retryWrites=true&w=majority
Environment=ROOT_URL=<%= @root_url %>
<%- if @mongo_replset -%>
Environment=MONGO_OPLOG_URL=mongodb+srv://<%= @mongo_host %>/local<%= authstring %><%= replstring %>&retryWrites=true&w=majority
<%- if @mongo_native_oplog -%>
Environment=USE_NATIVE_OPLOG=true
<%- end -%>
<%- end -%>
<%- if @speakeasy_http_header_name && @speakeasy_http_header_value -%>
Environment=SPEAKEASY_HTTP_HEADER_NAME=<%= @speakeasy_http_header_name %>
Environment=SPEAKEASY_HTTP_HEADER_VALUE=<%= @speakeasy_http_header_value %>
<%- end -%>

[Install]
WantedBy=multi-user.target
